# Configs/systemd/user/niri.service
[Unit]
Description=Niri Wayland Compositor

[Service]
ExecStart=/usr/bin/niri
Restart=on-failure

[Install]
WantedBy=default.target
